# crontab/start.sh
#!/bin/bash

open -a Terminal.app /Users/tuyuwang/Lagou/crawl.sh