# crontab/start.sh
#!/bin/bash

open -a Terminal.app F:/python/chonghcong/lagou/Lagou/crawl.sh