# crontab/update.sh
#!/bin/bash

open -a Terminal.app /Users/tuyuwang/Lagou/ip.sh