# crontab/update.sh
#!/bin/bash

open -a Terminal.app F:/python/chonghcong/lagou/Lagou/ip.sh